# crontab
*/5 * * * * root /bin/unrar_torrent.sh "/data/"

